# packages/automation-viavi/automation_viavi-1.0.tar.gz/automation_viavi-1.0/systemd/lab-automation.service
[Unit]
Description=Lab Automation Service
After=network.target

[Service]
Type=simple
Environment="API_TOKEN=NTIxNzMxMDUzODY5OtBkvhyICxJ4Dk8hN2vzOhs6slT4"
Environment="PAGE_ID=226939533"
ExecStart=/usr/local/bin/automation
Restart=on-failure
User=pi  

[Install]
WantedBy=multi-user.target