# packages/automation-viavi/automation_viavi-1.6-py3-none-any.whl/etc/systemd/system/lab-automation.timer
[Unit]
Description=Runs Lab Automation Service every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min
Unit=lab-automation.service

[Install]
WantedBy=timers.target
